# systemd/com.vovanr.onair.timer
[Unit]
Description=Onair Webcamera and Microphone Schedule

[Timer]
AccuracySec=1s
OnBootSec=5
OnUnitActiveSec=5

[Install]
WantedBy=timers.target
